# services/examples/demo-socket-echo/quadlet/bloom-demo-socket-echo.socket
[Unit]
Description=Bloom Demo Socket Echo — Socket activation listener

[Socket]
ListenStream=127.0.0.1:9090
Accept=no
Service=bloom-demo-socket-echo.service
SocketMode=0660

[Install]
WantedBy=sockets.target
