# bridges/whatsapp/quadlet/bloom-whatsapp.container
[Unit]
Description=Bloom WhatsApp Bridge (Baileys)
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/bloom-whatsapp:latest
ContainerName=bloom-whatsapp

# Connect to host network so it can reach bloom-channels on localhost:18800
Network=host

# Persist WhatsApp auth state across restarts
Volume=bloom-whatsapp-auth:/data/auth:Z

# Environment
Environment=BLOOM_CHANNELS_HOST=127.0.0.1
Environment=BLOOM_CHANNELS_PORT=18800
Environment=BLOOM_AUTH_DIR=/data/auth
Environment=NODE_ENV=production

# Resource limits (~50MB expected)
PodmanArgs=--memory=128m

# Health check
HealthCmd=node -e "process.exit(0)"
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s

# Security hardening
NoNewPrivileges=true
DropCapability=ALL

# Logging
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

[Install]
WantedBy=default.target
