# services/gateway/quadlet/bloom-gateway.container
[Unit]
Description=Bloom Gateway — unified web proxy (Cinny + Matrix + WebDAV)
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/bloom-gateway:latest
ContainerName=bloom-gateway

# Host networking — gateway must reach localhost-bound backends
Network=host

# Dynamic Caddyfile — regenerated when services are added/removed
Volume=%h/.config/bloom/Caddyfile:/etc/caddy/Caddyfile:ro

PodmanArgs=--memory=128m
PodmanArgs=--security-opt label=disable
HealthCmd=wget -qO- http://localhost:18810/ || exit 1
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s
HealthStartPeriod=30s
NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=default.target
